# gemma-bridge.service
[Unit]
Description=Gemma Bridge Monitor Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/USERNAME/PROJECTS_DIR/GemmaHelper/bridge_monitor.py
WorkingDirectory=/home/USERNAME/PROJECTS_DIR/GemmaHelper
Restart=on-failure
RestartSec=5
Type=simple

[Install]
WantedBy=default.target
